# Rpi/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.fish
# generated from rosbash/env-hooks/15.rosbash.fish.em

. "/home/pi/catkin_ws/src/src/ros/rosbash/rosfish"
