# Rpi/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.bash
# generated from rosbash/env-hooks/15.rosbash.bash.em

. "/home/pi/catkin_ws/src/src/ros/rosbash/rosbash"
