# Rpi/devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.tcsh
# generated from rosbash/env-hooks/15.rosbash.tcsh.em

. "/home/pi/catkin_ws/src/src/ros/rosbash/rostcsh"
