# Rpi/devel_isolated/rosbuild/etc/catkin/profile.d/10.rosbuild.sh
# generated from rosbuild/env-hooks/10.rosbuild.sh.em

# env variables in develspace
export ROS_ROOT="/home/pi/catkin_ws/src/src/ros/rosbuild"
